# rc.local
(cd /srv/blackboard && sudo -u hkim node app.js) &

